# config/securix-conf/etc/profile.d/sx-proxy.sh
#!/bin/bash

#: info: Securix GNU/Linux user/profile proxy setup
#: file: /etc/profile.d/sx-proxy.sh
#: descr: Put proxy settings for users here (http_proxy, ftp_proxy, ...)
#: author: Martin Cmelik (cm3l1k1) - securix.org, security-portal.cz
#: version: 20141123
#

